# ops/systemd/trade-heartbeat.timer
[Unit]
Description=Run trade heartbeat every 2 minutes

[Timer]
OnCalendar=*:0/2
AccuracySec=10s
Persistent=true
Unit=trade-heartbeat.service

[Install]
WantedBy=timers.target
