# ops/systemd/trade-heartbeat.timer
[Unit]
Description=Run trade heartbeat every 2 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=2min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
